# tools/internal/baoxiao/launchd/com.gengrowth.baoxiao-month-start.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gengrowth.baoxiao-month-start</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>@BAOXIAO_DIR@/launchd/run-month-start-on-first-day.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>@BAOXIAO_DIR@</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>@BAOXIAO_DIR@/logs/launchd-month-start.log</string>
    <key>StandardErrorPath</key>
    <string>@BAOXIAO_DIR@/logs/launchd-month-start.err</string>
    <key>RunAtLoad</key>
    <false/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>LC_ALL</key>
        <string>en_US.UTF-8</string>
        <key>PYTHONIOENCODING</key>
        <string>utf-8</string>
    </dict>
</dict>
</plist>
